# build-docker/debian/router-core.service
[Unit]
Description=GWRS Mini-Gateway Core
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/local/bin/router-core
Restart=on-failure
RestartSec=5
StandardOutput=append:/var/log/gwrs/core.log
StandardError=append:/var/log/gwrs/core.error.log
# Log rotation settings
LogRateLimitBurst=10000
LogsDirectory=gwrs
LogsDirectoryMode=0755
RuntimeMaxUse=50M
RuntimeMaxFileSize=20M
RuntimeMaxFiles=10
SystemMaxUse=200M
SystemMaxFileSize=50M
SystemMaxFiles=10

[Install]
WantedBy=multi-user.target